# conf/etc/s6/services/01-phpfpm/run
#!/usr/bin/with-contenv sh
service-hanlder "PHP-FPM" "php-fpm7 -F"
 